# system.d/passphrase.service
[Unit]
Description=Password Attribute Daemon

[Service]
ExecStart=%h/bin/passphrase.d

[Install]
WantedBy=default.target
